# AdminSide/sarima_api/sarima-api.service
[Unit]
Description=SARIMA Crime Forecast API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/alertdavao/AdminSide/sarima_api
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
